# papersome/rootfs/etc/s6-overlay/s6-rc.d/papersome-init/up
#!/bin/bash
exec /usr/local/bin/papersome-init.sh
