# papersome/rootfs/etc/s6-overlay/s6-rc.d/papersome-scheduler/run
#!/bin/bash
set -e

cd /var/www/html
echo "[papersome-scheduler] Starting Laravel scheduler..."
exec php artisan schedule:work
