# papersome/rootfs/etc/s6-overlay/s6-rc.d/papersome-queue/run
#!/bin/bash
set -e

cd /var/www/html
echo "[papersome-queue] Starting Laravel queue worker..."
exec php artisan queue:work --tries=1 --sleep=1
